# root/etc/s6-overlay/s6-rc.d/init-obsidian-config/up
/etc/s6-overlay/s6-rc.d/init-obsidian-config/run
